# build_system/deb_amd64/metadata/DEBIAN/postrm
#!/usr/bin/env bash

clear_symlinks()    {
    rm /usr/local/bin/usr_canet200_to_can_converter
}

clear_environment() {
    USERNAME="usr_canet200_converter"
    supervisorctl reload

    groupdel $USERNAME
    userdel $USERNAME
}

systemctl stop usr_canet200_to_can_converter
systemctl disable usr_canet200_to_can_converter

clear_symlinks

case "$1" in
     remove)
        clear_environment
        ;;
    *)
        ;;
esac
