# build_system/deb_amd64/metadata/DEBIAN/prerm
systemctl stop usr_canet200_to_can_converter
systemctl disable usr_canet200_to_can_converter
